# discord/feds-bot/feds-bot.openrc
#!/sbin/openrc-run

name="FEDS Discord Bot"
description="FEDS FRC Team 201 Discord Bot"

command="/usr/bin/node"
command_args="index.js"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
directory="/root/feds-central/discord/feds-bot"

depend() {
    need net
    after firewall
}

start_pre() {
    # Source environment variables
    if [ -f "${directory}/.env" ]; then
        export $(cat "${directory}/.env" | grep -v '^#' | xargs)
    fi
}
