# discord/feds-bot/feds-bot.service
[Unit]
Description=FEDS Discord Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/feds-central/discord/feds-bot
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=feds-bot

# Environment file
EnvironmentFile=/root/feds-central/discord/feds-bot/.env

[Install]
WantedBy=multi-user.target
